# scripts/linux/systemd/edgecoder-coordinator.service
[Unit]
Description=EdgeCoder Coordinator Runtime
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=edgecoder
Group=edgecoder
WorkingDirectory=__REPO_ROOT__
EnvironmentFile=/etc/edgecoder/coordinator.env
ExecStart=__NODE_BIN__ __REPO_ROOT__/dist/swarm/coordinator.js
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
