# scripts/macos/package-scripts/postinstall
#!/usr/bin/env bash

set -euo pipefail

SERVICE_LABEL="io.edgecoder.runtime"
PLIST_SRC="/opt/edgecoder/install/io.edgecoder.runtime.plist"
PLIST_DST="/Library/LaunchDaemons/${SERVICE_LABEL}.plist"
ENV_EXAMPLE="/etc/edgecoder/edgecoder.env.example"
ENV_FILE="/etc/edgecoder/edgecoder.env"
CONFIG_WIZARD="/opt/edgecoder/bin/edgecoder-configure.sh"
OLLAMA_INSTALLER="/opt/edgecoder/bin/edgecoder-install-ollama.sh"
LOG_DIR="/var/log/edgecoder"

mkdir -p "$LOG_DIR"
chmod 755 "$LOG_DIR"

if [[ ! -f "$ENV_FILE" && -f "$ENV_EXAMPLE" ]]; then
  cp "$ENV_EXAMPLE" "$ENV_FILE"
  chmod 600 "$ENV_FILE"
fi

# Package scripts run in a non-interactive PackageKit context, so we cannot
# reliably prompt here. Keep defaults and print explicit next-step command.
echo "EdgeCoder installed. Configure runtime with:"
echo "  sudo $CONFIG_WIZARD"

# Bootstrap Ollama as part of package install (best-effort, non-fatal).
if [[ -x "$OLLAMA_INSTALLER" ]]; then
  "$OLLAMA_INSTALLER" || true
fi

if [[ -f "$PLIST_SRC" ]]; then
  cp "$PLIST_SRC" "$PLIST_DST"
  chown root:wheel "$PLIST_DST"
  chmod 644 "$PLIST_DST"
fi

/bin/launchctl bootout system "$PLIST_DST" >/dev/null 2>&1 || true
/bin/launchctl bootstrap system "$PLIST_DST" || {
  echo "Warning: could not bootstrap ${SERVICE_LABEL} during install."
  echo "Run manually after fixing runtime env:"
  echo "  sudo launchctl bootstrap system $PLIST_DST"
  echo "  sudo launchctl kickstart -k system/${SERVICE_LABEL}"
  exit 0
}
/bin/launchctl enable "system/${SERVICE_LABEL}" || true
/bin/launchctl kickstart -k "system/${SERVICE_LABEL}" || true

exit 0
