# scripts/macos/payload/opt/edgecoder/install/io.edgecoder.runtime.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>io.edgecoder.runtime</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/opt/edgecoder/bin/edgecoder-runtime.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/opt/edgecoder/app</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/var/log/edgecoder/runtime.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/edgecoder/runtime.err.log</string>
  </dict>
</plist>
